# provisioning/ansible/roles/envcheck/files/envcheck.service
[Unit]
Description=Check Environments for ISUCON13 Contestants
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/opt/isucon-env-checker/warmup.sh
ExecStart=/opt/isucon-env-checker/run-isucon-env-checker.sh

[Install]
WantedBy=multi-user.target