# provisioning/ansible/roles/envcheck/files/envcheck.service
[Service]
ExecStart=/home/isucon/envcheck/envcheck boot
Type=oneshot
RemainAfterExit=yes
